# agent.win.ini
[program:collector]
command=%(here)s/../../embedded/python.exe ./agent/agent.py foreground --use-local-forwarder
directory=%(here)s/../../modules/datadog
environment=LANG=POSIX,PYTHONPATH=%(here)s/../../modules/datadog/lib
startsecs=5
startretries=3
stopasgroup=true
killasgroup=true
redirect_stderr=true
priority=998

[program:forwarder]
command=%(here)s/../../embedded/python.exe ./agent/ddagent.py --use_simple_http_client=1
directory=%(here)s/../../modules/datadog
environment=LANG=POSIX,PYTHONPATH=%(here)s/../modules/datadog/lib
startsecs=5
startretries=3
stopasgroup=true
killasgroup=true
redirect_stderr=true
priority=998